# build_files/root/usr/lib/systemd/system/setup-kargs.service
[Unit]
Description=Set kernel args on first boot
DefaultDependencies=no
After=systemd-remount-fs.service local-fs.target systemd-udev-settle.service
Before=plymouth-quit.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/plymouth display-message --text="Setting kernel arguments, please wait"; /usr/local/bin/reset-kargs.sh'
RemainAfterExit=yes
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
