# build_files/root/usr/lib/systemd/user/borshevik-app-manager-first-run.service
[Unit]
Description=Borshevik App Manager (first run)
After=graphical-session.target
Wants=graphical-session.target
ConditionPathExists=!%h/.local/state/borshevik/app-manager-first-run.done

[Service]
Type=simple
ExecStart=/usr/bin/borshevik-app-manager-first-run
Restart=no

[Install]
WantedBy=default.target
