# build_files/root/usr/lib/systemd/user/borshevik-kill-gnome-if-hung.service
[Unit]
Description=Kill GNOME Shell if it stops responding (D-Bus timeout)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/borshevik/borshevik-kill-gnome-if-hung.sh
TimeoutStartSec=60